# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
user=root
command=/bin/bash -c "/run_ssh.sh"

